# deploy/wodoo.service
[Unit]
Description=Wodoo - Inventory Location Manager
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/wodoo/backend
Environment="PATH=/opt/wodoo/venv/bin"
EnvironmentFile=/opt/wodoo/.env
ExecStart=/opt/wodoo/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/wodoo

[Install]
WantedBy=multi-user.target
